# stage-zbitx/02-zbitx-os/files/rootfs/etc/systemd/system/hostapd.service.d/10-zbitx-ap-channel.conf
# Match the AP channel to the wlan0 STA channel before hostapd starts.
# brcmfmac can't do ACS (channel=0), and AP+STA shares one channel — see
# zbitx-ap-channel and hostapd.conf.
[Unit]
After=uap0.service
Requires=uap0.service

[Service]
ExecStartPre=/usr/local/sbin/zbitx-ap-channel
