# stage-zbitx/02-zbitx-os/files/rootfs/etc/systemd/system/hostapd.service.d/restart.conf
# Auto-restart hostapd so the `zbitx` AP self-heals — e.g. if hostapd loses the
# uap0 interface or races uap0.service at boot. The zBitx v2 WiFi front panel
# connects over this AP, so a dead AP means a dead front panel.
#
# Restart=on-failure (not =always) so a deliberate stop — e.g. the AP-manager
# GUI's "Stop AP" — is honoured. RestartSec=5 keeps retries ~5s apart, well
# under systemd's default start-limit burst, so hostapd keeps retrying instead
# of giving up after the interface eventually appears.
[Service]
Restart=on-failure
RestartSec=5
