# stage-zbitx/02-zbitx-os/files/rootfs/etc/systemd/system/zbitx-ap-follow-channel.service
[Unit]
Description=Reconcile the zbitx AP channel with the wlan0 STA channel
# Pointless unless the AP is actually running.
After=hostapd.service
ConditionPathExists=/sys/class/net/wlan0

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/zbitx-ap-reconcile
