# layer/files/etc/systemd/system/zbitx-expand-rootfs.service
[Unit]
Description=Expand root filesystem to fill the storage device (first boot)
Documentation=https://github.com/arunderwood/zbitxv2-image
# rpi-image-gen's image-rpios builds a fixed-size image (root partition
# sized to the rootfs contents); nothing else grows it, so without this
# the card boots ~100% full. Run once, before the display manager so the
# graphical session has space to write ~/.Xauthority etc.
ConditionVirtualization=!container
ConditionPathExists=!/var/lib/zbitx/rootfs-expanded
After=local-fs.target
Before=lightdm.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/zbitx-expand-rootfs
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
