# layer/files/etc/systemd/system/getty@tty1.service.d/autologin.conf
# Auto-login the `pi` user on tty1 — mirrors what `raspi-config` writes
# when you pick "Console with autologin" (B2).
# Combined with /home/pi/.bash_profile's startx invocation, this gives
# us a console boot that lands directly in the sBitx GTK UI without
# the weight of a full display manager.
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin pi --noclear %I $TERM
